# sync-diag.service
[Unit]
Requires=mnt-heapdumps.mount
After=mnt-heapdumps.mount

[Service]
ExecStart=/opt/rchain-testnet-node/sync-diag
